# launchd/com.fhirmini.ai.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Template: __ROOT__ is replaced with the repo path by `scripts/fhirmini install`. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fhirmini.ai</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__ROOT__/scripts/run-ai.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>__ROOT__/logs/ai.out</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/logs/ai.err</string>
</dict>
</plist>
